# upstartDMap.conf
description "Configures DMap (git pull from GNRS-DOS, sets this server's IP addr using bash_setIP) and starts DMap"
author "Colleen Rock"
start on startDMap
stop on shutdown

script 
   cd /root/GNRS-DOS/
   exec /root/GNRS-DOS/start_DMap_server.sh
end script
